# scripts/ai.watson.knowledge-intake-poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.watson.knowledge-intake-poll</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <!-- Replace with path to poll-channel.js in your skill installation -->
    <string>/PATH/TO/skills/x-bookmark-triage/scripts/poll-channel.js</string>
  </array>
  <key>StartInterval</key>
  <integer>120</integer>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/tmp/knowledge-intake-poll.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/knowledge-intake-poll.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <!-- Copy required vars from your gateway plist or set explicitly -->
    <key>DISCORD_BOT_TOKEN</key>
    <string>YOUR_BOT_TOKEN</string>
    <key>ANTHROPIC_DEFAULT_KEY</key>
    <string>YOUR_ANTHROPIC_KEY</string>
    <key>KNOWLEDGE_INTAKE_CHANNEL_ID</key>
    <string>YOUR_CHANNEL_ID</string>
    <key>OPENCLAW_WORKSPACE</key>
    <string>/PATH/TO/workspace</string>
  </dict>
</dict>
</plist>
